# docker/artwork_export/crontab
00 0 * * * /var/www/html/manage/export_cron.sh
